# .config/systemd/user/jack-carla-rack.service.d/override.conf
[Service]
Environment='CARLA_CONFIG_FILE=/home/jaco/Documents/Carla/laptop_setup.carxp'
ExecStart=
ExecStart=/usr/bin/pw-jack carla-rack -n $CARLA_CONFIG_FILE
